# systemd/backup-system-daily.timer.example
[Unit]
Description=Daily System Backup Timer (Example)
Documentation=file:///opt/backup-system/docs/README.md
Requires=backup-system@system.service

[Timer]
Unit=backup-system@system.service
# Run daily at 10:00 AM
OnCalendar=*-*-* 10:00:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
